# scheduler/crontab
0 5 * * * python3 /scheduler/send_task.py delete_old_tasks
